# agent/systemd/vpn-agent.service
[Unit]
Description=ICD360S VPN management agent
Documentation=https://github.com/ICD360S-e-V/vpn
After=network-online.target wg-quick@wg0.service AdGuardHome.service
Wants=network-online.target
Requires=wg-quick@wg0.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/vpn-agent serve
Restart=on-failure
RestartSec=5

# Hardening — see systemd.exec(5) and systemd-analyze security vpn-agent
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictNamespaces=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
# SystemCallFilter is intentionally NOT set in M4.4: modernc.org/sqlite
# (and the Go runtime it links against) calls some syscall outside the
# @system-service whitelist on AlmaLinux 9's systemd (255). Without
# SystemCallErrorNumber=, the kernel kills the process with SIGSYS
# instead of returning EPERM, and the journal does not say which call
# was blocked. Until we identify and explicitly allow the offending
# call, the filter is disabled. CAP_NET_BIND_SERVICE + CAP_NET_ADMIN
# remain the primary authorisation boundary; ProtectSystem,
# ProtectKernelModules, NoNewPrivileges, MemoryDenyWriteExecute, and
# RestrictNamespaces still apply. TODO M5: re-enable with the right
# whitelist after running with SystemCallErrorNumber=EPERM to identify
# the missing syscall.

# State directories the agent is allowed to write to.
# M2: /etc/wireguard is now RW because the agent mutates wg0.conf when
# peers are added or removed via the API.
ReadWritePaths=/etc/vpn-agent /var/lib/vpn-agent /etc/wireguard

# Capabilities — keep minimal
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_ADMIN

# Logging
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=4096
LimitNPROC=512

[Install]
WantedBy=multi-user.target
